# scripts/fetch_indices.cron
0 17 * * 1-5 python /app/fetch_indices.py >> /var/log/cron.log 2>&1
